# x728.service
[Unit]
Description=Monitors x728 and shuts down Pi

[Service]
Type=simple
ExecStart=python /opt/kismet-x728/x728.py
